# images/gateway/registry.container
# Podman quadlet (system): the fleet OCI registry. Serves both images
# (seafront app + seafront-os) to the boxes over the isolated backbone.
# Installed to /etc/containers/systemd/ by gateway-setup.sh.
[Unit]
Description=Fleet OCI registry
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/library/registry:2
# gateway's fleet-LAN IP; boxes pull from here
PublishPort=10.10.0.69:5000:5000
# persistent image store on the writable /var partition
Volume=/var/lib/seafront-registry:/var/lib/registry:Z

[Service]
Restart=always

[Install]
WantedBy=multi-user.target default.target
