# images/kinoite/files/etc/containers/systemd/seafront.container
# Podman quadlet: runs the seafront app container as a system service.
# `systemctl restart seafront` works because quadlet generates seafront.service.
[Unit]
Description=seafront microscope control
After=network-online.target
Wants=network-online.target

[Container]
Image=10.10.0.69:5000/seafront:stable
PublishPort=8000:8000
# Persistent state on the writable host partition (survives OS upgrades).
# config.json is pushed here by the gateway (scripts/push-config.sh); images land here too.
Volume=/home/pharmbio/seafront:/root/seafront:z
# SEAFRONT_HOME must equal $HOME/seafront (=/root/seafront) — see the app Containerfile.
Environment=SEAFRONT_HOME=/root/seafront
# Profile to run from config.json. Flip to the real profile once hardware is wired.
Environment=SEAFRONT_MICROSCOPE=mocroscope

# ---------------------------------------------------------------------------
# USB PASSTHROUGH — validated on squid1 (2026-07-20). Two host-side prerequisites
# are handled by seafront-hardware-prep.service (baked in the OS image), NOT here:
#   - cameras need the SELinux `container_use_devices` boolean (raw usbfs/libusb);
#   - the Teensy `--device` is dropped by this quadlet at boot (see below) and
#     re-attached by that unit after udev settles.
# Robust-to-hotplug pattern: mount the whole usb tree + allow the device-cgroup
# classes, rather than binding a single node that vanishes on replug.
#   189 = usb (libusb cameras), 166 = ttyACM (Teensy CDC-ACM serial)
# ---------------------------------------------------------------------------
Volume=/dev/bus/usb:/dev/bus/usb
# NOTE the quotes: each cgroup-rule value contains spaces, and quadlet splits
# PodmanArgs on whitespace — unquoted, the unit parses into garbage.
PodmanArgs=--device-cgroup-rule='c 189:* rmw' --device-cgroup-rule='c 166:* rmw'
# The Teensy microcontroller is a USB CDC-ACM serial device; seafront finds/controls it
# via pyserial, which needs the /dev/ttyACM* NODE inside the container (the /dev/bus/usb
# mount above only gives raw libusb access, not the serial node). Pass it through. The
# leading '-' makes each optional, so a box with no Teensy — e.g. running the mock profile
# — still starts. The cgroup rule above already grants access to char-major-166 (ttyACM).
#
# GOTCHA (why seafront-hardware-prep.service exists): quadlet's optional-device '-' check
# runs in the systemd GENERATOR at early boot, before udev creates /dev/ttyACM0, so it
# evaluates "absent" and drops --device every boot. The prep unit daemon-reloads after
# udev settles (regenerating this unit WITH --device) and restarts seafront. SELinux is
# NOT involved for --device nodes (podman relabels them; they open fine with the boolean off).
AddDevice=-/dev/ttyACM0
AddDevice=-/dev/ttyACM1

[Service]
Restart=always
RestartSec=5
# Refresh from the gateway registry on start if reachable; keep the local image when offline.
ExecStartPre=-/usr/bin/podman pull 10.10.0.69:5000/seafront:stable

[Install]
WantedBy=multi-user.target default.target
